# launchd/com.aaron.pilot-tune-reviews.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aaron.pilot-tune-reviews</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/aaron/development/pilot/scripts/tune-reviews.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/aaron/Documents/Claude/outputs/pilot-tune-reviews-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aaron/Documents/Claude/outputs/pilot-tune-reviews-launchd.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/aaron/.npm-global/bin</string>
        <key>HOME</key>
        <string>/Users/aaron</string>
    </dict>
</dict>
</plist>
