# launchd/com.aaron.pilot-triage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aaron.pilot-triage</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/aaron/Documents/Scripts/lift-triage.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <!-- Sunday 22:30 -->
        <dict>
            <key>Weekday</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- Tuesday 22:30 -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <!-- Thursday 22:30 -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/aaron/Documents/Claude/outputs/pilot-triage-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aaron/Documents/Claude/outputs/pilot-triage-launchd.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/aaron/.npm-global/bin</string>
        <key>HOME</key>
        <string>/Users/aaron</string>
    </dict>
</dict>
</plist>
